# gtfs-update-tbilisi.service
[Unit]
Description=Download latest Tbilisi GTFS static feed
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot

# Download to a temp file then atomically move into place so
# gtfs-realtime never sees a partial zip.
ExecStart=/bin/sh -c '\
    curl -fL --retry 3 --retry-delay 30 \
        -o /var/lib/gtfs-realtime-tbilisi/gtfs.zip.tmp \
        https://jbb.ghsq.de/gtfs/ge-tbilisi.gtfs.zip \
    && mv /var/lib/gtfs-realtime-tbilisi/gtfs.zip.tmp \
          /var/lib/gtfs-realtime-tbilisi/gtfs.zip'

# Restart the RT service so it picks up the new feed.
# '+' prefix runs this line as root, outside the sandbox.
ExecStartPost=+/bin/systemctl restart gtfs-realtime-tbilisi.service

# Run as the same dynamic user that owns the state directory
DynamicUser=yes
StateDirectory=gtfs-realtime-tbilisi
WorkingDirectory=/var/lib/gtfs-realtime-tbilisi

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gtfs-update-tbilisi

# Filesystem sandboxing
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadOnlyPaths=/
ReadWritePaths=/var/lib/gtfs-realtime-tbilisi

# Network — only needs outbound HTTPS
RestrictAddressFamilies=AF_INET AF_INET6

# Namespace & privilege isolation
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
PrivateDevices=yes
MemoryDenyWriteExecute=yes

# System call filtering
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @swap @reboot @obsolete

# Misc hardening
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077
ProcSubset=pid
ProtectProc=invisible
RemoveIPC=yes
